# prototype/com.wtyuyuyu.secondmind.prototype.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.wtyuyuyu.secondmind.prototype</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/Users/wtyuyuyu/Library/Application Support/SecondMindSite/server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/wtyuyuyu/Library/Application Support/SecondMindSite</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PYTHONDONTWRITEBYTECODE</key>
      <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/wtyuyuyu/Library/Application Support/SecondMindSite/logs/launchd.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/wtyuyuyu/Library/Application Support/SecondMindSite/logs/launchd.stderr.log</string>
  </dict>
</plist>
