# prototype/com.wtyuyuyu.secondmind.latest-focus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.wtyuyuyu.secondmind.latest-focus</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/Users/wtyuyuyu/Library/Application Support/SecondMindSite/update_latest_focus.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/wtyuyuyu/Library/Application Support/SecondMindSite</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PYTHONDONTWRITEBYTECODE</key>
      <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/wtyuyuyu/Library/Application Support/SecondMindSite/logs/latest-focus.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/wtyuyuyu/Library/Application Support/SecondMindSite/logs/latest-focus.stderr.log</string>
  </dict>
</plist>
